# api.conf
[program:api]
command=./run_server.bash
directory=/home/entrant/yandex-backend-school-task/src
user=entrant
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true